# files/etc/profile.d/orientdb.sh
#!/usr/bin/env sh
export PATH=$PATH:/opt/orientdb/bin
export ORIENTDB_HOME=/opt/orientdb
